# haproxy/haproxy.ini
[supervisord]
loglevel=debug
nodaemon=true

[program:confd]
command=/usr/local/bin/confd -node %(ENV_ETCD_NODE)s

[program:haproxy]
command=/usr/sbin/haproxy -db -f /etc/haproxy/haproxy.cfg
