# nginx/nginx.ini
[supervisord]
loglevel=debug
nodaemon=true

[program:confd]
command=/usr/local/bin/confd -node %(ENV_ETCD_NODE)s

[program:nginx]
command=/usr/sbin/nginx
